# contrib/openrc/monero-wallet-rpc
#!/sbin/openrc-run

description="Monero Wallet RPC"
command="/var/www/underground.pm/contrib/monero/run.sh"

depend() {
    need network
}

start() {
    ebegin "Starting Monero Wallet RPC"
    start-stop-daemon --start --exec $command --background --user root
    eend $?
}

stop() {
    ebegin "Stopping Monero Wallet RPC"
    start-stop-daemon --stop --exec $command
    eend $?
}

restart() {
	stop
	sleep 1
	start
}
