# contrib/openrc/monero-test-wallet-rpc
#!/sbin/openrc-run

description="Monero Test Wallet RPC"
command="monero_test_wallet_rpc_run"

depend() {
    need net
}

start() {
    ebegin "Starting Monero Test Wallet RPC"
    start-stop-daemon --start --exec $command --background --user root
    eend $?
}

stop() {
    ebegin "Stopping Monero Test Wallet RPC"
    start-stop-daemon --stop --exec $command
    eend $?
}

restart() {
	stop
	sleep 1
	start
}
