# contrib/openrc/underground
#!/sbin/openrc-run

description="Privacy hosting"
name=$RC_SVCNAME

command="/usr/bin/underground"
pidfile="/run/$name.pid"

depend() {
    need net
    need libvirtd
    need postgresql
    need monerod
}

start() {
    ebegin "Starting $name"
    start-stop-daemon --background --make-pidfile --start --pidfile "$pidfile" --exec "$command"
    eend $?
}

stop() {
    ebegin "Stopping $name"
    start-stop-daemon --stop --pidfile "$pidfile"
    eend $?
}

restart() {
    stop
    sleep 1
    start
}
